# lib/canonical/buildd/debian/launchpad-buildd.cron.daily
#!/bin/sh

find /var/log/launchpad-buildd/ -mtime +2 -print0 | xargs -r -0 rm
